# plugins/idna/supervisor/conf.d/idna.conf
[program:idna]
command=uv run %(ENV_HOME)s/.roxabi/idna/idna_server.py
directory=%(ENV_HOME)s/.roxabi/idna
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:%(ENV_PATH)s"
autostart=false
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/supervisor/logs/idna.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/supervisor/logs/idna_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
