# plugins/forge/supervisor/conf.d/forge.conf
[program:forge]
command=%(ENV_HOME)s/.roxabi/forge/run.sh
directory=%(ENV_HOME)s/.roxabi/forge
environment=HOME="%(ENV_HOME)s",PATH="%(ENV_HOME)s/.local/bin:/usr/local/bin:/usr/bin:/bin",FORGE_DIR="%(ENV_HOME)s/.roxabi/forge"
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/.local/state/supervisor/logs/forge.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/.local/state/supervisor/logs/forge_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
